# ugw/target/linux/x86/base-files-openwrt/etc/uci-defaults/dhcp
#!/bin/sh

uci set dhcp.dnsmasq.interface='br-lan'
uci del_list dhcp.dnsmasq.listenaddress='127.0.0.1'
uci del_list dhcp.dnsmasq.listenaddress='192.168.1.1'
uci commit dhcp

exit 0
