# ugw/target/linux/x86/base-files-openwrt/etc/init.d/inetd
#!/bin/sh /etc/rc.common

START=24
USE_PROCD=1

INETD_CONFIG="/tmp/inetd.conf"
INETD_PROG="/usr/sbin/inetd"
INETD_PID="/var/run/inetd.pid"

add_inetd_service()
{
	local service_name="$1"
	local socket_type="$2"
	local protocol="$3"
	local nowait="${4/0/nowait}"; nowait="${nowait/1/wait}"
	local username="$5"
	local program="$6"
	shift 6
	local program_with_args="$@"

	if [ -n "$service_name" -a -n "$socket_type" -a -n "$protocol" -a -n "$nowait" -a -n "$username" -a -n "$program" ]; then
		[ -f "$program" ] && {
			echo -e "$service_name\t$socket_type\t$protocol\t$nowait\t$username\t$service_name\tugw$service_name\t-\t-\t$program\t$program_with_args" >> $INETD_CONFIG
		}
	fi
}

configure_inetd()
{
	> $INETD_CONFIG

	# add_inetd_service <service name> <socket type> <protocol> <wait 1/nowait 0> <user> <program path> <program with args>
	add_inetd_service tftp   dgram  udp 1 root /usr/bin/tftpd     /usr/bin/tftpd -c /tmp
	add_inetd_service telnet stream tcp 1 root /usr/sbin/telnetd  /usr/sbin/telnetd -l /bin/login  -w 1
	add_inetd_service ftp    stream tcp 0 root /usr/sbin/vsftpd   /bin/nice -n 19 /usr/sbin/vsftpd
}

start_service() {
	configure_inetd
	procd_open_instance
	procd_set_param command $INETD_PROG $INETD_CONFIG -f
	procd_set_param respawn
	procd_close_instance
}

reload_service() {
	local pid=$(cat $INETD_PID 2>/dev/null)
	[ -z "$pid" ] && {
		start_service
		return 0;
	}

	configure_inetd
	kill -HUP $pid
}

stop_service() {
	service_stop $INETD_PROG
	[ -f "$INETD_CONFIG" ] && {
		rm -f "$INETD_CONFIG"
	}
}

