# ugw/target/linux/lantiq/base-files/etc/uci-defaults/ethernet
#!/bin/sh
#
# Copyright (C) 2011-2012 OpenWrt.org
#

[ -e /etc/config/ethernet ] && exit 0

add_eth_lan_port() {
	local port=$1
	local enable=$2
	local speed=$3
	local duplex=$4
	local td=$5
	uci batch <<EOF
set ethernet.${port}=port
set ethernet.${port}.enable='${enable}'
set ethernet.${port}.speed='${speed}'
set ethernet.${port}.duplex='${duplex}'
set ethernet.${port}.td='${td}'
EOF
}

add_eth_wan_port() {
	local port=$1
	local enable=$2
	local speed=$3
	local duplex=$4
	local wan=$5
	uci batch <<EOF
set ethernet.${port}=port
set ethernet.${port}.enable='${enable}'
set ethernet.${port}.speed='${speed}'
set ethernet.${port}.duplex='${duplex}'
set ethernet.${port}.wan='${wan}'
EOF
}

add_trafficdesc() {
	local trafficdesc=$1
	local mbr=$2
	local mbs=$3
	local rate=$4
	uci batch <<EOF
set ethernet.${trafficdesc}=trafficdesc
set ethernet.${trafficdesc}.mbr='${mbr}'
set ethernet.${trafficdesc}.mbs='${mbs}'
set ethernet.${trafficdesc}.rate='${rate}'
EOF
}


. /lib/functions/uci-defaults.sh
. /lib/lantiq.sh

touch /etc/config/ethernet

board=$(lantiq_board_name)

case "$board" in
EASY350*)
	add_eth_lan_port "eth0_1" "1" "1000" "full"
	add_eth_lan_port "eth0_2" "1" "1000" "full"
	add_eth_lan_port "eth0_3" "1" "1000" "full"
	add_eth_lan_port "eth0_4" "1" "1000" "full" "shapingconf"

	add_eth_wan_port "eth1" "1" "1000" "full" "1"
	;;

EASY330\ VDSL\ BOND)
	add_eth_lan_port "eth0_1" "1" "1000" "full"
	add_eth_lan_port "eth0_2" "1" "1000" "full"
	add_eth_lan_port "eth0_3" "1" "1000" "full"
	add_eth_lan_port "eth0_4" "1" "1000" "full" "shapingconf"

	add_eth_wan_port "eth1" "1" "1000" "full" "1" 
	
	add_trafficdesc "shapingconf" "20000" "2000" "enabled"

	;;

EASY330\ ETH\ RT)
	add_eth_lan_port "eth0_1" "1" "1000" "full"
	add_eth_lan_port "eth0_2" "1" "1000" "full"
	add_eth_lan_port "eth0_3" "1" "1000" "full"
	add_eth_lan_port "eth0_4" "1" "1000" "full" "shapingconf"

	add_eth_wan_port "eth1" "1" "1000" "full" "1" 
	
	add_trafficdesc "shapingconf" "20000" "2000" "enabled"

	;;

EASY300\ AC1200)
	add_eth_lan_port "eth0_1" "1" "1000" "full"
	add_eth_lan_port "eth0_2" "1" "1000" "full"

	add_eth_wan_port "eth1" "1" "1000" "full" "1"

	add_trafficdesc "shapingconf" "20000" "2000" "enabled"
	;;
esac

uci commit ethernet

exit 0
