# ugw/feeds_opensrc/packages/lighttpd/files/lighttpd
#!/bin/sh /etc/rc.common

/etc/preinit.d/scapi_utils.sh
USE_PROCD=1
START=10

start_service() {
        procd_open_instance
        procd_set_param command lighttpd -D -f /etc/lighttpd/lighttpd.conf
        procd_set_param respawn
        procd_close_instance
}

