# ugw/feeds_opensrc/packages/samba3/files/uci-defaults/samba
#!/bin/sh

grep -q ^ugw: /etc/passwd || {
	adduser ugw -u 201 -G root -h /mnt/usb -g ugw -D -H
	(echo ugw;echo ugw;echo ugw;)|passwd ugw
}

grep -q ^ugw: /etc/samba/smbpasswd || {
	smbpasswd ugw ugw
}

exit 0

