# ugw/feeds_opensrc/packages/netifd/files/etc/init.d/network_stop
#!/bin/sh /etc/rc.common

STOP=02

stop() 
{
	br2684ctl -K 2>/dev/null
	ubus call service delete '{ "name": "netifd" }'
}

