# ugw/target/linux/lantiq/base-files-openwrt/etc/hotplug.d/iface/10-internet-led
#!/bin/sh
# Hotplug script to auto configure Internet LED based on primary WAN
#

. /lib/functions.sh
. /lib/functions/uci-defaults.sh


[ "$ACTION" = "ifup" -a -n "$INTERFACE" ] && {
	interface=`head -2 /proc/net/route|tail -1|awk '{ print $1 }'`
	value=`head -2 /proc/net/route|tail -1|awk '{ print $2 }'`
	config_load network
	config_get if_name $INTERFACE ifname
	[ "$if_name" = "$interface" ] && [ "$value" = "00000000" ] && {
		uci -q get system.led_internet && {
			device=`uci get system.led_internet.dev`
			[ "$device" != "$interface" ] && {
				uci set system.led_internet.dev="$interface"
				uci set system.led_internet.conn="$INTERFACE"
				uci commit system
				/etc/init.d/led reload
			} || {
				intf=`uci get system.led_internet.conn`
				[ "$intf" != "$INTERFACE" ] && uci set system.led_internet.conn="$INTERFACE" && uci commit system
				echo 255 > /sys/class/leds/internet/brightness
			}
		} || {
			ucidef_set_led_netdev "internet" "internet" "internet" "$interface"
			uci set system.led_internet.conn="$INTERFACE"
			uci commit system
			echo 255 > /sys/class/leds/internet/brightness
		}
	}
}

[ "$ACTION" = "ifdown" -a -n "$INTERFACE" ] && {
	value=`head -2 /proc/net/route|tail -1|awk '{ print $2 }'`
	[ "$value" != "00000000" ] && echo 0 > /sys/class/leds/internet/brightness
}
