# ugw/target/linux/lantiq/base-files-openwrt/etc/hotplug.d/iface/30-tunnels
[ "$ACTION" == "ifup" ] && {
	proto=$(uci get network.$INTERFACE.proto)
	[ "$proto" == "gre" ] && {
		tunipaddr=$(uci get network.$INTERFACE.tunipaddr)
		tunnetmask=$(uci get network.$INTERFACE.tunnetmask)
		ip addr add $tunipaddr dev gre-$INTERFACE
		ip route add $tunnetmask dev gre-$INTERFACE

		# Add WAN-to-LAN FORWARD Rule for IPoGRE pkts
		iptables -I zone_wan_forward -s $tunnetmask -j zone_lan_dest_ACCEPT
	}
	[ "$proto" == "dslite" ] && {
		# Add WAN-to-LAN FORWARD Rule for DS-Lite pkts
		iptables -I zone_wan_forward -s 192.0.0.1 -j zone_lan_dest_ACCEPT
	}
	[ "$proto" == "6rd" ] && {
		ip6prefix=$(uci get network.$INTERFACE.ip6prefix)
		ip6prefixlen=$(uci get network.$INTERFACE.ip6prefixlen)
		# Add WAN-to-LAN FORWARD Rule for 6rd pkts
		ip6tables -I zone_wan_forward -s $ip6prefix/$ip6prefixlen -j zone_lan_dest_ACCEPT
	}
}

[ "$ACTION" == "ifdown" ] && {
	proto=$(uci get network.$INTERFACE.proto)
	[ "$proto" == "gre" ] && { 
		tunnetmask=$(uci get network.$INTERFACE.tunnetmask)

		# Add WAN-to-LAN FORWARD Rule for IPoGRE pkts
		iptables -D zone_wan_forward -s $tunnetmask -j zone_lan_dest_ACCEPT
	}
	[ "$proto" == "dslite" ] && {
		# Add WAN-to-LAN FORWARD Rule for DS-Lite pkt
		iptables -D zone_wan_forward -s 192.0.0.1 -j zone_lan_dest_ACCEPT
	}
	[ "$proto" == "6rd" ] && {
		ip6prefix=$(uci get network.$INTERFACE.ip6prefix)
		ip6prefixlen=$(uci get network.$INTERFACE.ip6prefixlen)
		# Add WAN-to-LAN FORWARD Rule for 6rd pkts
		ip6tables -D zone_wan_forward -s $ip6prefix/$ip6prefixlen -j zone_lan_dest_ACCEPT
	}
}
