# ugw/target/linux/lantiq/base-files-ugw-newframework/etc/init.d/free_cache
#!/bin/sh /etc/rc.common

START=90
start() {
	/bin/free_caches.sh &
}

stop() {
killall -9 free_caches.sh
}

restart() {
stop
sleep 2
start
}

