# ugw/target/linux/lantiq/base-files-ugw-newframework/etc/init.d/inetd
#!/bin/sh /etc/rc.common

START=04

INETD_CONFIG="/etc/inetd.conf"

chk_update_cfg()
{
	[ -f "$1" ] && {
		echo "$2" >> $INETD_CONFIG
	}
}

start()
{
	local inetdpid;

	echo -n > /etc/inetd.conf
	chk_update_cfg /usr/bin/tftpd     "tftp   dgram  udp6 wait   root tftp   ifxftp  - - /usr/bin/tftpd    tftpd -c /ramdisk/tftp_upload"
	chk_update_cfg /usr/sbin/telnetd  "telnet stream tcp6 wait   root telnet ifxrcon - - /usr/sbin/telnetd telnetd -l /bin/login -w 1"

	mkdir -p /var/run/vsftpd
	chk_update_cfg /usr/sbin/vsftpd   "ftp    stream tcp nowait root vsftpd vsftpd  - - /bin/nice         nice -n 19 /usr/sbin/vsftpd"

	if [ -f /var/run/inetd.pid ]; then
		#The server is already running and so just send HUP singal
		inetdpid=`cat /var/run/inetd.pid`
		if [ -d /proc/$inetdpid ];then
			kill -HUP $inetdpid
		else
			rm -f /var/run/inetd.pid
			/usr/sbin/inetd $INETD_CONFIG 2>/dev/null
		fi
	else
		#The server is not running and so start it
		/usr/sbin/inetd $INETD_CONFIG 2> /dev/null
	fi
	dmesg | grep GCC | cut -f2 -d')' | cut -c3- > /tmp/toolchain-ver
}

stop()
{
	killall inetd
}

