# ugw/target/linux/lantiq/base-files-ugw-newframework/usr/bin/led_status
#!/bin/sh /etc/rc.common

######################################################################
#
# Function  : init script, telnetd_startup
#
# Creater   : xiangda.zhang@feixun.com.cn, 2016-05-04
#
# Copyright : Shanghai Phicomm
#
######################################################################

START=99

boot() {
  /usr/bin/led_status &
}

start() {
  /usr/bin/led_status &
}

stop() {
  killall -9 led_status
}

restart() {
  stop
  start
}
