# ugw/target/linux/x86/base-files-ugw-newframework/etc/init.d/version.sh
#!/bin/sh /etc/rc.common
#START=59

start() {
	/usr/sbin/version.sh -r&
}
