# ugw/feeds_opensrc/packages/mwan3/files/etc/init.d/mwan3
#!/bin/sh /etc/rc.common

START=03

start()
{
	echo "Inlcuding modules needed for mwan3 ..." > /dev/console
	insmod nfnetlink

	insmod ip_set
	insmod ip_set_bitmap_ip
	insmod ip_set_bitmap_ipmac
	insmod ip_set_bitmap_port
	isnmod ip_set_hash_ip
	insmod ip_set_hash_ip
	insmod ip_set_hash_ipport
	insmod ip_set_hash_ipportip
	insmod ip_set_hash_ipportnet
	insmod ip_set_hash_net
	insmod ip_set_hash_netiface
	insmod ip_set_hash_netport
	insmod ip_set_list_set

	insmod xt_set
	insmod xt_CT
	insmod xt_HL
	insmod xt_NFQUEUE
	insmod xt_comment
	insmod xt_connbytes
	insmod xt_connlimit
	insmod xt_ecn
	insmod xt_hl
	insmod xt_layer7
	insmod xt_recent
	insmod xt_statistic
}
