# ugw/feeds_opensrc/packages/mountd/files/mountd.init
#!/bin/sh /etc/rc.common
START=80

start()
{
  insmod /lib/modules/*/autofs4.ko
  sync
  echo 1 > /proc/sys/vm/drop_caches
  sleep 1
  [ -d /mnt/usb/ ] || {
	[ -h /mnt/usb ] || ln -sf /tmp/usb /mnt/usb
	mkdir -p /tmp/usb
  }
  /sbin/mountd
}

stop()
{
  killall mountd
  rmmod ahci
  umount /ramdisk/usb/.run/mountd || true
  rmmod autofs4
}
