# ugw/target/linux/lantiq/base-files-ugw/etc/init.d/powered_up
#!/bin/sh /etc/rc.common

START=24

start() {
	if [ "$CONFIG_TARGET_LTQCPE_PLATFORM_AR9_VB" = "1" ]; then
		echo default-on > /sys/class/leds/pwr:green/trigger
	fi
}
