# package/network/services/ipwatchd/files/ipwatchd.init
#!/bin/sh /etc/rc.common

START=80

start() {
	conf="/etc/ipwatchd.conf"

	if [ -f "$conf" ]; then
		ipwatchd -c $conf
	else
		logger -t ipwatchd "Invalid configuration $conf"
	fi
}

stop() {
	killall ipwatchd
}