# ugw/target/linux/intel-xway/xrx500/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop

#ttyLTQ0::askfirst:/bin/ash --login
#ttyLTQ1::askfirst:/bin/ash --login

::respawn:/sbin/getty 115200 ttyLTQ0

#ttyS0::askfirst:/bin/ash --login
