# ugw/target/linux/lantiq/base-files-openwrt/etc/uci-defaults/ppa
#!/bin/sh
#

[ -e /etc/config/ppa ] && exit 0

add_ppa()
{
	local iftype=$1
	local interfaces="$2"
	uci batch <<EOF
set ppa.${iftype}='ppa'
set ppa.${iftype}.type='${iftype}'
set ppa.${iftype}.ifname='${interfaces}'
EOF
}

. /lib/functions/uci-defaults.sh
. /lib/lantiq.sh
touch /etc/config/ppa

board=$(lantiq_board_name)

case "$board" in
EASY330\ VDSL\ BOND)
	add_ppa lan "eth0"
	add_ppa wan "eth1 ptm0 nas0"
	;;

EASY350*)
	add_ppa lan "eth0_1 eth0_2 eth0_3 eth0_4"
	add_ppa wan "eth1 ptm0 nas0"
	;;

EASY300\ AC1200)
	add_ppa lan "eth0"
	add_ppa wan "eth1 ptm0 nas0"
	;;
Phicomm*)
	add_ppa lan "eth0_1 eth0_2 eth0_3 eth0_4"
	add_ppa wan "eth1"
	;;
esac

uci commit ppa

exit 0
