# ugw/target/linux/lantiq/base-files-openwrt/etc/uci-defaults/uhttpd
#!/bin/sh

## By default uhttpd https_redirection is on :(
uci set uhttpd.main.redirect_https='0'
uci commit uhttpd

exit 0
